# crontab
# Запуск каждый день в 01:00 ночи
0 1 * * * python /app/main.py >> /var/log/cron.log 2>&1
